# etc/udev/rules.d/70-persistent-monitor.rules
#!/bin/bash
#KERNEL=="card0", SUBSYSTEM=="drm", ACTION=="change", RUN+="/bin/su jzb -c '/home/jzb/bin/monitor_hotplug.sh'"
